# docker/supervisord.conf
; Supervisord configuration for coding-services container
; Manages all MCP servers and web services

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ===========================================
; MCP Servers (HTTP/SSE)
; ===========================================

[program:semantic-analysis]
command=node /coding/integrations/mcp-server-semantic-analysis/dist/sse-server.js
directory=/coding/integrations/mcp-server-semantic-analysis
environment=
    CODING_ROOT="/coding",
    SEMANTIC_ANALYSIS_PORT="%(ENV_SEMANTIC_ANALYSIS_PORT)s",
    QDRANT_URL="%(ENV_QDRANT_URL)s",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    GOOGLE_AI_API_KEY="%(ENV_GOOGLE_AI_API_KEY)s"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/semantic-analysis.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/semantic-analysis-error.log
stderr_logfile_maxbytes=10MB
priority=100

[program:browser-access]
command=node /coding/integrations/browser-access/dist/sse-server.js
directory=/coding/integrations/browser-access
environment=
    BROWSER_ACCESS_PORT="%(ENV_BROWSER_ACCESS_PORT)s",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/browser-access.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/browser-access-error.log
stderr_logfile_maxbytes=10MB
priority=100

[program:constraint-monitor]
command=node /coding/integrations/mcp-constraint-monitor/src/sse-server.js
directory=/coding/integrations/mcp-constraint-monitor
environment=
    CONSTRAINT_MONITOR_PORT="%(ENV_CONSTRAINT_MONITOR_PORT)s",
    QDRANT_URL="%(ENV_QDRANT_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s",
    HTTP_PROXY="",
    HTTPS_PROXY="",
    http_proxy="",
    https_proxy=""
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/constraint-monitor.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/constraint-monitor-error.log
stderr_logfile_maxbytes=10MB
priority=100

[program:code-graph-rag]
command=/coding/integrations/code-graph-rag/.venv/bin/python -m codebase_rag.mcp.sse_server
directory=/coding/integrations/code-graph-rag
environment=
    CODE_GRAPH_RAG_PORT="%(ENV_CODE_GRAPH_RAG_PORT)s",
    MEMGRAPH_HOST="%(ENV_MEMGRAPH_HOST)s",
    MEMGRAPH_PORT="%(ENV_MEMGRAPH_PORT)s",
    TARGET_REPO_PATH="/workspace",
    HTTP_PROXY="",
    HTTPS_PROXY="",
    http_proxy="",
    https_proxy=""
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/code-graph-rag.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/code-graph-rag-error.log
stderr_logfile_maxbytes=10MB
priority=100

; ===========================================
; Web Services
; ===========================================

[program:vkb-server]
command=node /coding/lib/vkb-server/express-server.js
directory=/coding/lib/vkb-server
environment=
    VKB_PORT="%(ENV_VKB_PORT)s",
    QDRANT_URL="%(ENV_QDRANT_URL)s",
    HTTP_PROXY="",
    HTTPS_PROXY="",
    http_proxy="",
    https_proxy=""
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/vkb-server.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/vkb-server-error.log
stderr_logfile_maxbytes=10MB
priority=50

[program:health-dashboard]
command=node /coding/integrations/system-health-dashboard/server.js
directory=/coding/integrations/system-health-dashboard
environment=
    PORT="%(ENV_HEALTH_DASHBOARD_PORT)s",
    WS_PORT="%(ENV_HEALTH_DASHBOARD_WS_PORT)s"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/health-dashboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/health-dashboard-error.log
stderr_logfile_maxbytes=10MB
priority=50

[program:health-dashboard-frontend]
command=node /coding/integrations/system-health-dashboard/static-server.js
directory=/coding/integrations/system-health-dashboard
environment=HEALTH_DASHBOARD_PORT="%(ENV_HEALTH_DASHBOARD_PORT)s"
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=5
stdout_logfile=/var/log/supervisor/health-dashboard-frontend.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/health-dashboard-frontend-error.log
stderr_logfile_maxbytes=5MB
priority=50

; ===========================================
; Health Monitoring
; ===========================================

[program:health-verifier]
command=node /coding/scripts/health-verifier.js start
directory=/coding
environment=CODING_REPO="/coding"
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/health-verifier.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/health-verifier-error.log
stderr_logfile_maxbytes=10MB
priority=200

; ===========================================
; Process Groups
; ===========================================

[group:mcp-servers]
programs=semantic-analysis,browser-access,constraint-monitor,code-graph-rag
priority=100

[group:web-services]
programs=vkb-server,health-dashboard,health-dashboard-frontend
priority=50

[group:monitoring]
programs=health-verifier
priority=200
